# dublin-mirror/systemd-services/plo-w4p.service
[Unit]
Description=W4P Remote Control Flask (potlimitomaha.xyz)
After=network.target docker.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/plo-w4p
EnvironmentFile=/opt/plo-w4p/.env
ExecStart=/opt/plo-w4p/venv/bin/gunicorn -c gunicorn.conf.py app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plo-w4p

# Resource limits
MemoryMax=1G
MemoryHigh=768M

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/plo-w4p /opt/plo-engine/static /opt/plo-equity /tmp
ProtectHome=true
NoNewPrivileges=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
